# bkwww/bkwww.cron.daily
#!/usr/bin/bash
/usr/local/bin/bkwww.sh
